# src/skcapstone/data/systemd/skcapstone@.service
# SKCapstone multi-agent systemd template unit.
#
# Each named agent runs as a separate service instance:
#
#   systemctl --user enable skcapstone@opus skcapstone@jarvis
#   systemctl --user start  skcapstone@opus skcapstone@jarvis
#
# The instance name (%i) is passed to --agent so the daemon picks up
# ~/.skcapstone/agents/%i/ as its home directory and auto-assigns a port.
#
# Stop a single agent:
#   systemctl --user stop skcapstone@jarvis

[Unit]
Description=SKCapstone Sovereign Agent — %i
Documentation=https://github.com/smilinTux/skcapstone
After=network-online.target ollama.service syncthing.service
Wants=network-online.target
# Multiple instances may coexist; no conflicts between them.
# Crash-loop guard: give up after StartLimitBurst failures inside the
# interval instead of restarting forever. This is the .41 outage failure
# mode (a failed preflight looping every 10s). With the exponential backoff
# below, 6 restarts land at ~610s cumulative, well inside the 1800s window,
# so a persistently failing daemon stops and stays failed within one window.
StartLimitIntervalSec=1800
StartLimitBurst=6
# Page out-of-band when this agent enters the failed state (crash-loop
# give-up or hard failure) instead of failing silently.
OnFailure=skcapstone-alert@%i.service

[Service]
Type=notify
ExecStart=%h/.local/bin/skcapstone daemon start --agent %i --foreground
ExecStop=%h/.local/bin/skcapstone daemon stop --agent %i
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
# Exponential restart backoff: 10s, 20s, 40s ... capped at 5 min, so a
# persistently failing daemon does not hot-loop every 10s (systemd >= 254).
RestartSteps=5
RestartMaxDelaySec=300
# Memory caps encode the .41 host fix in the unit (was hand-applied host
# state only). Normal agent RSS is ~230 MB; the 4G hard cap matches the
# legacy single-agent unit and the .41 host value, giving ~17x headroom
# while still stopping a runaway before it OOM-thrashes the box. MemoryHigh
# applies reclaim pressure first so a slow leak throttles before the kill.
MemoryHigh=3G
MemoryMax=4G
# Watchdog: daemon must call sd_notify("WATCHDOG=1") at least every 5 minutes
WatchdogSec=300
# Pass instance name as env var so imported modules that read SKCAPSTONE_AGENT
# resolve the correct per-agent home even before the CLI flag is processed.
Environment=PYTHONUNBUFFERED=1
Environment=OLLAMA_KEEP_ALIVE=5m
Environment=SKCAPSTONE_AGENT=%i
# Journal logging — logs appear under skcapstone@<instance>
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skcapstone@%i

# Security hardening (same as single-agent unit)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.skcapstone %h/.skmemory %h/.local %h/.config %h/.capauth %h/.cloud9 %h/.skcomms %h/.skchat
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=default.target
