# src/skcapstone/data/systemd/skcomms-heartbeat.service
[Unit]
Description=SKComms Heartbeat Emission
Documentation=https://github.com/smilinTux/skcomms
After=network-online.target

[Service]
Type=oneshot
ExecStart=skcomms heartbeat --no-emit
ExecStart=skcomms heartbeat
Nice=19

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.skcapstone %h/.skcomms
PrivateTmp=true

Environment=PYTHONUNBUFFERED=1
