# src/skcapstone/data/systemd/skcomms-queue-drain.timer
[Unit]
Description=SKComms Queue Drain Timer — periodic message retry
Documentation=https://github.com/smilinTux/skcomms

[Timer]
OnBootSec=60s
OnUnitActiveSec=120s
RandomizedDelaySec=15s
Persistent=true

[Install]
WantedBy=timers.target
