# src/skcapstone/data/systemd/skcomms-heartbeat.timer
[Unit]
Description=SKComms Heartbeat Timer — periodic liveness announcements
Documentation=https://github.com/smilinTux/skcomms

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
RandomizedDelaySec=10s
Persistent=true

[Install]
WantedBy=timers.target
