# systemd/skcapstone.service
[Unit]
Description=SKCapstone Sovereign Agent Daemon
Documentation=https://github.com/smilinTux/skcapstone
After=network-online.target ollama.service syncthing.service
Wants=network-online.target
# Crash-loop guard: give up after StartLimitBurst failures inside the
# interval instead of restarting forever (the .41 outage failure mode).
StartLimitIntervalSec=1800
StartLimitBurst=6
# Page out-of-band when the daemon enters the failed state.
OnFailure=skcapstone-alert@skcapstone.service

[Service]
Type=simple
ExecStart=%h/.skenv/bin/skcapstone daemon start --foreground
ExecStop=%h/.skenv/bin/skcapstone daemon stop
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
# Exponential restart backoff: 10s, 20s, 40s ... capped at 5 min, so a
# persistently failing daemon does not hot-loop every 10s (systemd >= 254).
RestartSteps=5
RestartMaxDelaySec=300
# Cap memory to prevent OOM from large model loading. MemoryHigh applies
# reclaim pressure first so a slow leak throttles before the hard kill.
MemoryHigh=3G
MemoryMax=4G
# Keep Ollama models warm for 5 minutes between requests
Environment=PYTHONUNBUFFERED=1
Environment=OLLAMA_KEEP_ALIVE=5m
Environment=SKAGENT=lumina
Environment=SKCAPSTONE_AGENT=lumina
Environment=SKMEMORY_AGENT=lumina
# Journal logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skcapstone

# Security hardening (relaxed — ProtectHome=read-only breaks if any
# ReadWritePaths dir is missing on the host)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
