# systemd/skcomm-heartbeat.service
[Unit]
Description=SKComm Heartbeat Emission
Documentation=https://github.com/smilinTux/skcomm
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.skenv/bin/skcomm heartbeat --no-emit
ExecStart=%h/.skenv/bin/skcomm heartbeat
Nice=19

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.skcapstone %h/.skcomm
PrivateTmp=true

Environment=PYTHONUNBUFFERED=1
Environment=SKAGENT=lumina
Environment=SKCAPSTONE_AGENT=lumina
Environment=PATH=%h/.skenv/bin:/usr/local/bin:/usr/bin:/bin
